# bin/system/iot-display.service
[Unit]
Description=iot-display
Documentation=https://github.com/RaederDev/iot-display
After=network.target

[Service]
Type=simple
User=display
ExecStart=/opt/iot-display/bin/launch-server.sh
WorkingDirectory=/opt/iot-display/
Restart=no

[Install]
WantedBy=multi-user.target
